# deploy-ubuntu/telegram-bot-api.service
[Unit]
Description=Telegram Bot FastAPI (Uvicorn)
After=network.target

[Service]
Type=simple
User=nopolicybot
Group=nopolicybot
WorkingDirectory=/opt/telegram-political-monitor-bot
Environment=PATH=/opt/telegram-political-monitor-bot/venv/bin
ExecStart=/opt/telegram-political-monitor-bot/venv/bin/python /opt/telegram-political-monitor-bot/run_api.py
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/opt/telegram-political-monitor-bot/data
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
