# deploy-ubuntu/telegram-bot.service
[Unit]
Description=Telegram Political Monitor Bot
After=network.target

[Service]
Type=notify
NotifyAccess=main
User=root
Group=root
WorkingDirectory=/opt/telegram-political-monitor-bot
Environment=PATH=/opt/telegram-political-monitor-bot/venv/bin
EnvironmentFile=-/opt/telegram-political-monitor-bot/.env
ExecStart=/opt/telegram-political-monitor-bot/venv/bin/python /opt/telegram-political-monitor-bot/bot.py
Restart=always
RestartSec=10
KillSignal=SIGINT
# Длинные LLM-запросы / aiogram shutdown — иначе SIGKILL через 20s и потеря очереди апдейтов
TimeoutStopSec=90
WatchdogSec=300
# LLM + граф + кэши; 512M часто упирается в OOM после часов работы
MemoryMax=768M
MemorySwapMax=256M

[Install]
WantedBy=multi-user.target
