# deploy-ubuntu/telegram-bot-admin.service
[Unit]
Description=Telegram Bot Admin (Gunicorn)
After=network.target

[Service]
Type=simple
User=nopolicybot
Group=nopolicybot
WorkingDirectory=/opt/telegram-political-monitor-bot
Environment=PATH=/opt/telegram-political-monitor-bot/venv/bin
ExecStart=/bin/bash -lc 'CPU="$(nproc 2>/dev/null || echo 1)"; W="$((2 * CPU + 1))"; if [ "$W" -lt 3 ]; then W=3; fi; exec /opt/telegram-political-monitor-bot/venv/bin/gunicorn -w "$W" -b 127.0.0.1:5000 --timeout 30 --keep-alive 5 admin_app:app'
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/opt/telegram-political-monitor-bot/data
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
